# com.sandykoufax.fantasy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sandykoufax.fantasy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sethsteiman/fantasy-baseball-agent/run_daily.sh</string>
    </array>

    <!-- Fire at 7:00 AM local time (PT = 10 AM ET — before any MLB game starts) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the machine was asleep at 7 AM, run as soon as it wakes up -->
    <key>StandardOutPath</key>
    <string>/Users/sethsteiman/fantasy-baseball-agent/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sethsteiman/fantasy-baseball-agent/logs/launchd.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/sethsteiman/fantasy-baseball-agent</string>

    <!-- Do NOT run immediately when loaded — only at scheduled time -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
